# Linux in Action/Глава 5. АВтоматизированное администрирование/site-backup.service
[Unit]
Description=Backup Apache website

[Service]
Type=simple
ExecStart=/home/user/site-backup.sh

[Install]
WantedBy=multi-user.target