# Linux in Action/Глава 5. АВтоматизированное администрирование/site-backup.timer
[Unit]
Description =Backup Apache website timer

[Timer]
OnCAlendar=*-*-* 5:51:00
Unit=site-backup.sercive

[Install]
WantedBy=multi-user.target